# archive/legacy/launchd-backups/20260403-001543/com.wan.gtd.sync.user.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wan.gtd.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/wan/workspace/gtd-tasks/scripts/sync_agent_mac.py</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/Users/wan/workspace/gtd-tasks/logs/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wan/workspace/gtd-tasks/logs/sync-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>GTD_API_BASE_URL</key>
        <string>https://gtd.5666.net</string>
    </dict>
</dict>
</plist>
